# inspur-bios-code-mgmt.service
[Unit]
Description=Update BIOS

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/sbin/inspur-bios-updater-sh
SyslogIdentifier=image-bios
