# inspur-bios-code-mgmt.service.in
[Unit]
Description=Update BIOS

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/sbin/inspur-bios-code-mgmt
SyslogIdentifier=image-bios
